# deploy/glance-web.service
[Unit]
Description=glance web (tailnet-only summarizer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zen
WorkingDirectory=/home/zen/dev/glance
# Bind to localhost; expose to the tailnet via `tailscale serve` on :8443
# (TLS-terminated by Tailscale). See README for the one-liner.
Environment=GLANCE_HOST=127.0.0.1
Environment=GLANCE_PORT=8765
ExecStart=/home/zen/.local/bin/uv run glance-web
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
